# systemd/forgejo-stack.service
[Unit]
Description=Forgejo (Git hosting + runner) stack (rootless podman-compose)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=PODMAN_USERNS=keep-id
WorkingDirectory=%h/podman-compose/forgejo
ExecStart=/usr/bin/podman-compose up -d
ExecStop=/usr/bin/podman-compose down
RemainAfterExit=yes

[Install]
WantedBy=default.target

