# systemd/openqa-stack.service
[Unit]
Description=OpenQA stack (rootless podman-compose)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/podman-compose/openqa/prod
ExecStart=/usr/bin/podman-compose up -d
ExecStop=/usr/bin/podman-compose down
RemainAfterExit=yes

[Install]
WantedBy=default.target

